# bb-epaper.service
[Unit]
Description=BB-Terminal e-paper companion service (renders frames for XIAO ESP32-S3 panel)
After=network-online.target
Wants=network-online.target
# Wait for the BB-Terminal OpenBB API on :6900 to be reachable; not strictly required at boot
# (cache falls back to placeholder values on first miss), but keeps the first frame meaningful.

[Service]
Type=simple
User=jason
Group=jason
WorkingDirectory=/usr/local/bb-epaper
ExecStart=/usr/local/bb-epaper/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8081 --log-level info
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
